# conf/get-mlb-schedule.cron
# m h  dom mon dow user command
00 00,12 * * * www-data /usr/local/sbin/get-mlb-schedule.sh
